# roles/user-x/files/systemd/wm.service
[Unit]
Description=Window manager

[Service]
ExecStart=/usr/bin/notion
